# systemd/i4ops-backend.service
[Unit]
Description=i4ops Dashboard Backend API
Documentation=https://github.com/your-org/i4ops-dashboard
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=i4ops
Group=i4ops
WorkingDirectory=/home/i4ops/i4ops-dashboard/backend
Environment=NODE_ENV=production
EnvironmentFile=/home/i4ops/i4ops-dashboard/.env
ExecStart=/usr/bin/node dist/index.js
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/i4ops/i4ops-dashboard
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=i4ops-backend

[Install]
WantedBy=multi-user.target 